# dotfiles/config/systemd/user/autostart.target.wants/tint2.service
[Unit]
Description=Lightweight panel

[Service]
ExecStart=/usr/bin/tint2

Restart=always
RestartSec=5s

[Install]
WantedBy=autostart.target