# dotfiles/config/systemd/user/dropbox.service
[Unit]
Description=Dropbox autostart

[Service]
ExecStart=/usr/bin/dropbox start