# dotfiles/config/systemd/user/nm-applet.service
[Unit]
Description=Systray applet to manage Network Manager

[Service]
ExecStart=/usr/bin/nm-applet
