# dotfiles/config/systemd/user/redshift-gtk.service
[Unit]
Description=Color temperature adjustment tool

[Service]
ExecStart=/usr/bin/redshift-gtk -l 44:10
