# dotfiles/config/systemd/user/tint2.service
[Unit]
Description=Lightweight panel

[Service]
ExecStart=/usr/bin/tint2
